# quadlets/prometheus.container.j2
[Container]
AutoUpdate=registry
ContainerName=prometheus
Group=0
HealthCmd=wget -nv -t1 --spider http://localhost:9090/metrics
HealthOnFailure=kill
Image={{ prometheus_image }}
LogDriver=journald
Network=services.network
Timezone=local
User=0
Volume={{ prometheus_data_volume }}:/prometheus:Z

[Service]
Restart=always

[Install]
WantedBy=default.target
