# quadlets/mirror/archlinux/web.container.j2
[Container]
AutoUpdate=registry
ContainerName=mirror-archlinux-web
HostName=archlinux.mirror
Image={{ mirror_web_image }}
LogDriver=journald
Network=proxy.network
Timezone=local
Volume={{ mirror_archlinux_volume }}:/mirror:z,ro

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
